# docker/s6-rc.d/hermes-webui/run
#!/command/with-contenv sh
# shellcheck shell=sh
# Hermes WebUI (nesquena) on loopback — fronted by the control-plane proxy.

if [ "${TERMINAL_HOME_MODE:-}" = "real" ]; then
    export HOME="/opt/data"
else
    export HOME="${HOME:-${HERMES_DATA_DIR:-/opt/data}/.hermes/home}"
fi
export HERMES_DATA_DIR="${HERMES_DATA_DIR:-/opt/data}"
export HERMES_HOME="${HERMES_HOME:-${HERMES_DATA_DIR}/.hermes}"
cd "$HOME"

# shellcheck disable=SC1091
. /opt/hermes/.venv/bin/activate

export PYTHONPATH="/app:${PYTHONPATH:-}"
export HERMES_WEBUI_AGENT_DIR="${HERMES_WEBUI_AGENT_DIR:-/opt/hermes}"
export HERMES_WEBUI_STATE_DIR="${HERMES_WEBUI_STATE_DIR:-${HERMES_DATA_DIR}/webui}"
export HERMES_WORKSPACE_DIR="${HERMES_WORKSPACE_DIR:-${HERMES_DATA_DIR}/workspace}"
export HERMES_CONFIG_PATH="${HERMES_CONFIG_PATH:-${HERMES_HOME}/config.yaml}"
export HERMES_WEBUI_HOST="${CONTROL_PLANE_INTERNAL_WEBUI_HOST:-127.0.0.1}"
export HERMES_WEBUI_PORT="${CONTROL_PLANE_INTERNAL_WEBUI_PORT:-8788}"
export CONTROL_PLANE_RUNTIME=s6

exec s6-setuidgid hermes python /app/vendor/hermes-webui/server.py
